# ser2net.init
#!/bin/sh /etc/rc.common
# Ser2net daemon

START=99
STOP=99

START=10
STOP=15

start() {
        ser2net
}

stop(){
        killall ser2net
}
